# systemd/find_menu.service
[Unit]
Description=Instagram Menu Image Crawler Service
After=network.target

[Service]
Type=simple
# FIXME: Change YOUR_USER to your actual Linux username
# NOTE: Ensure 'node' is in your PATH. If systemd fails to find 'node', add:
# Environment=PATH=/usr/bin:/usr/local/bin:/home/YOUR_USER/.nvm/versions/node/vXX.XX.XX/bin
WorkingDirectory=/home/YOUR_USER/find_menu_images
ExecStart=/home/YOUR_USER/find_menu_images/run_cron.sh --today-only --only-new --check-time-window --extract-text
User=YOUR_USER

[Install]
WantedBy=multi-user.target
